# target/linux/ramips/base-files/etc/init.d/bootcount
#!/bin/sh /etc/rc.common

START=99

start() {
	. /lib/functions.sh

	case $(board_name) in
	sk-wb8)
		fw_setenv bootcount 0
		;;
	ea8100)
		mtd resetbc s_env || true
		;;
	esac

}
